# dockers/web/supervisord-mongodb.conf
[program:mongodb]
command=/usr/bin/mongod
numprocs=1
autostart=true
autorestart=true