# dockers/guacamole/supervisor/guacamole.sv.conf
[program:guacamole]
command=guacd -f
startretries=3
autorestart=true
redirect_stderr=true