# dockers/base/supervisord-vncserver.conf
[program:vncserver]
command=/start-vncserver.sh
numprocs=1
autostart=true
autorestart=true