# dockers/office/supervisord-crontab.conf
[program:crontab]
command=/start-crontab.sh
numprocs=1
autostart=true
autorestart=true