# office/crontab
*/1 * * * * root /usr/bin/calibredb add /opt/calibre/import --library-path /opt/calibre/library >> /var/log/cron.log 2>&1